# templates/home-mike.automount.j2
[Unit]
Description=Automount Home Mike
Requires=network-online.target
After=network-online.service

[Automount]
Where=/home/mike

[Install]
WantedBy=multi-user.target
