# templates/home-mike.mount.j2
[Unit]
Description=Home Mike

[Mount]
What=files.maballou.com:/srv/home/mike
Where=/home/mike
Type=nfs

[Install]
WantedBy=multi-user.target

