# templates/home-mike.automount.j2
[Unit]
Description=Automount Home Mike

[Automount]
Where=/home/mike

[Install]
WantedBy=multi-user.target

