# templates/mnt-share.automount.j2
[Unit]
Description=Automount Share

[Automount]
Where=/mnt/share

[Install]
WantedBy=multi-user.target

