# roles/fedora-workstation/templates/mnt-share.automount.j2
[Unit]
Description=Automount Share
Requires=network-online.target
After=network-online.service

[Automount]
Where=/mnt/share

[Install]
WantedBy=multi-user.target
